# run/systemd/generator/rpcbind.service.d/50-rpcbind-$portmap.conf
# Automatically generated by systemd-insserv-generator

[Unit]
Wants=rpcbind.target
Before=rpcbind.target
