# etc/xdg/autostart/blueproximity-autostart.desktop
[Desktop Entry]
Exec=blueproximity
Name=blueproximity
Terminal=false
Type=Application
StartupNotify=false
